# sway-session-save.timer
[Unit]
Description=Sway session auto-save (every 5 minutes)

[Timer]
OnActiveSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
